# silos/systemd/.config/systemd/user/dns-refresh.timer
[Unit]
Description=Refresh DigitalOcean DNS (Hourly Timer)

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
