# silos/systemd/.config/systemd/user/dns-refresh.service
[Unit]
Description=Refresh DNS Service

[Service]
ExecStart=/home/hq1/bin/.secret/update_dns.sh
