# silos/systemd/.config/systemd/user/forex-alert.timer
[Unit]
Description=Check Forex alerts (Every 15 minutes)

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
